# edr/EDR/systemd/edr-agent-test.service
[Unit]
Description=EDR Agent (test environment)
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/lcz/EDR
ExecStart=/home/lcz/EDR/edr-agent --config configs/agent_test.json
Restart=always
RestartSec=0
StartLimitIntervalSec=0

# Hardening
NoNewPrivileges=true
LockPersonality=true
RestrictRealtime=true
MemoryDenyWriteExecute=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ReadWritePaths=/home/lcz/edr_test
CapabilityBoundingSet=CAP_KILL CAP_DAC_OVERRIDE CAP_NET_ADMIN CAP_SETUID CAP_SETGID CAP_DAC_READ_SEARCH CAP_SYS_ADMIN CAP_BPF
AmbientCapabilities=CAP_KILL CAP_DAC_OVERRIDE CAP_NET_ADMIN CAP_SYS_ADMIN CAP_BPF

CPUQuota=50%
MemoryMax=512M

[Install]
WantedBy=multi-user.target
